# external/firefly_ec200s/21-ec200s.rules
SUBSYSTEM=="usb-serial", ACTION=="add", KERNEL=="ttyUSB[1-2]", GOTO="dongle_usb_rules"
GOTO="dongle_rules_end"

LABEL="dongle_usb_rules"

ENV{ID_VENDOR_ID}=="2c7c", ENV{ID_MODEL_ID}=="6002", GOTO="quectel_dialing_start"

LABEL="quectel_dialing_start"

RUN+="/usr/bin/ec200s.sh"

GOTO="dongle_rules_end"

LABEL="dongle_rules_end"
